# deploy/artothek-web.service
[Unit]
Description=Artothek web app
After=network.target

[Service]
Type=simple
User=artothek
Environment=PORT=8080
ExecStart=/opt/artothek-ui/app/.venv/bin/gunicorn \
    --chdir /opt/artothek-ui/app \
    --workers 2 \
    --bind 0.0.0.0:${PORT} \
    --access-logfile /var/log/artothek/access.log \
    --error-logfile  /var/log/artothek/error.log \
    app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
