# deploy/artothek-scraper.service
[Unit]
Description=Artothek scraper
After=network.target

[Service]
Type=oneshot
User=artothek
WorkingDirectory=/opt/artothek-ui/scraper
ExecStart=/opt/artothek-ui/scraper/.venv/bin/python scraper.py
StandardOutput=journal
StandardError=journal
